# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

; ── FastAPI / uvicorn (priority 10 = starts first) ───────────────────────────
[program:uvicorn]
command=python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --log-level info
directory=/app/backend
environment=PYTHONPATH="/app/backend"
autostart=true
autorestart=true
startretries=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── nginx (priority 20 = starts after uvicorn) ───────────────────────────────
; nginx may briefly return 502 for /api/* until uvicorn finishes binding (~1s).
; Static SPA files are served immediately regardless.
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
